# testPmt/testdata/supervisor.d/test1.conf
[program:tes1:]
directory=/home/www/server/test1
environment=LD_LIBRARY_PATH="/home/www/server/test1/lib",JUNO_TIME="189092893"
command=/home/www/server/test1/bin/test1 --host=0.0.0.1 --config=/home/www/.config/juno-agent/wsd-projecta-msg-srv-chatroom-tim-go/config/config-pre.toml
user=this is a change
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/home/www/server/test1/stdout.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stdout_capture_maxbytes=100MB
stderr_logfile=/home/www/server/test1/stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stopsignal=QUIT

1231m1:1231-123
