# testPmt/testdata/supervisor.d/test2.conf
[program:test4]
directory=/home/www/server/test3
environment=LD_LIBRARY_PATH="/home/www/server/test2/lib",JUNO_TIME="1909890"
command=/home/www/server/test2/bin/test2 --host=127.0.0.1 --config=./config/config-pre.toml
user=www
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/home/www/server/test2/stdout.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stdout_capture_maxbytes=100MB
stderr_logfile=/home/www/server/test2/stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
stopsignal=QUIT
